# MouseDriver.service
[Unit]
Description=Userspace Driver for T16 Gaming Mouse

[Service]
ExecStart=/usr/bin/mouseDriver


[Install]
WantedBy=default.target
